# .claude/skills/llm-wiki/references/discover-watchlist-launchd.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  /llm-wiki discover-watchlist 用 launchd plist テンプレ（Phase 3g 同梱）

  このファイルは **例（テンプレ）** です。`launchctl load` は利用者の手動操作です（自動配置はしません）。
  refresh-tier-a-launchd.plist.example / refresh-watchlist-launchd.plist.example と同型。
  違いは (1) ProgramArguments が `/llm-wiki discover-watchlist --no-prompt`、
  (2) 起動時刻が **refresh-tier-a（03:00）/ refresh-watchlist（03:30）と stagger**（3 系統）。

  使い方:
    1. 下記の <<REPLACE_REPO_PATH>> を本リポジトリの実 absolute path に置換
       例:
         REPO=/Users/<you>/Projects/personal-wiki-for-claude-code
         sed -i '' "s#<<REPLACE_REPO_PATH>>#${REPO}#g" \
           "${REPO}/.claude/skills/llm-wiki/references/discover-watchlist-launchd.plist.example"
    2. ~/Library/LaunchAgents/ に **コピー**して plist を配置（例ファイルそのものを load しない）
         cp .claude/skills/llm-wiki/references/discover-watchlist-launchd.plist.example \
           ~/Library/LaunchAgents/local.llm-wiki-discover-watchlist.plist
    3. launchctl load ~/Library/LaunchAgents/local.llm-wiki-discover-watchlist.plist
    4. 動作確認のため一度だけ手動起動:
         launchctl start local.llm-wiki-discover-watchlist
       以降は StartCalendarInterval（既定 04:00 ローカル）で自動起動。
    5. 解除するときは:
         launchctl unload ~/Library/LaunchAgents/local.llm-wiki-discover-watchlist.plist
         rm  ~/Library/LaunchAgents/local.llm-wiki-discover-watchlist.plist

  StandardOutPath / StandardErrorPath の既定は vault 直下の隠しログファイル。
  vault Git の意図しないコミットを避けるため、vault `.gitignore` に次の 2 行を追加してください:
    .discover-watchlist-out.log
    .discover-watchlist-err.log

  -- 各設定キーの意図 --

  ProgramArguments:
    `claude --print '/llm-wiki discover-watchlist --no-prompt' --allowedTools=Read,Write,Edit,Bash,Grep,Glob,WebFetch`
    の固定起動。
    - `--no-prompt`: フィード巡回 + stage-1 フィルタ + pending_feed_discoveries append のみ実行。
      AskUserQuestion 不発火・stage-2 LLM 不実行（cron での LLM API コストゼロ）。
      stage-2 判定と ingest 承認は対話セッションで `/llm-wiki discover-watchlist` を実行して行う。
    - `--allowedTools=` は **`=` 形式必須**。Claude CLI は variadic 引数を取るため
      スペース区切り（`--allowedTools <list>`）にすると後続の prompt 位置引数を
      食い込む既知の罠がある（spike で確認・memory: claude-code-cli-allowedtools-variadic）。

  EnvironmentVariables (HOME / PATH / USER / LOGNAME / SHELL の 5 変数 + NOTION_API_KEY):
    launchd は sanitized env でプロセスを起動するため、最低限以下の 5 変数を明示しないと
    macOS keychain auth が解決できず `Not logged in` で fail する
    （spike で確認・memory: claude-code-keychain-auth-env-vars）。
    前提検証（非対話起動・WebFetch 権限・vault 書き込み・git commit）は Phase 3a で
    refresh-tier-a により verified 済。mode I は mode G/W と同一コードパス（curl + WebFetch）のため
    3g では再 spike 不要。

    **NOTION_API_KEY（Phase 4 追記）**: feed_registry の `notion-medium-db` ソース巡回で
    minitools `discover-notion-medium` を Bash 実行する際に必須（未設定だと CLI が exit 1 で
    当該 Notion source skip・他フィードは続行＝clean failure）。**Chrome / Playwright 系 env は不要**
    （cron `--no-prompt` は stage-2 / ingest を実行せず `scrape-medium` Playwright を発火しない＝D3。
    Playwright を使う Medium content fetch〔4a〕は対話セッションでのみ走る）。
    minitools は `uv run --directory <minitools_path>` で起動するため PATH に `uv`
    （通常 ~/.local/bin）が含まれることを確認（下記 PATH に既定で含む）。

  WorkingDirectory:
    本リポジトリ直下を CWD にする必要がある（`.llm-wiki.json` を本リポジトリ直下から相対解決するため）。

  StartCalendarInterval:
    既定 **04:00 ローカル**（深夜帯）。
    **3 系統 stagger 構成**:
      - refresh-tier-a:    03:00（mode F・Tier A 再取得）
      - refresh-watchlist: 03:30（mode W・Tier B watchlist 再取得）
      - discover-watchlist: **04:00**（mode I・Tier B フィード新着発見）
    同一 vault 直下 `.llm-wiki.lock`（atomic 取得・排他制御）を争うため、3 系統が同時刻起動だと
    2 系統が毎回 lock 競合 skip し合う。30〜60 分ずらすことで先発の処理が完了して lock を解放した後に
    後発が走る。別時刻にしたい場合も 3 系統が重ならない Hour / Minute にすること。

  RunAtLoad:
    false（load 時の即時起動はしない・想定外の二重実行回避）。

  ANTHROPIC_API_KEY:
    通常は keychain 経由で解決するため plist に埋めない。
    上記 5 変数を満たしても auth で fail する場合のみ、最後の手段として追加。

-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.llm-wiki-discover-watchlist</string>

    <key>ProgramArguments</key>
    <array>
        <string>claude</string>
        <string>--print</string>
        <string>/llm-wiki discover-watchlist --no-prompt</string>
        <string>--allowedTools=Read,Write,Edit,Bash,Grep,Glob,WebFetch</string>
    </array>

    <key>WorkingDirectory</key>
    <string><<REPLACE_REPO_PATH>></string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/<<REPLACE_USER>></string>
        <key>USER</key>
        <string><<REPLACE_USER>></string>
        <key>LOGNAME</key>
        <string><<REPLACE_USER>></string>
        <key>SHELL</key>
        <string>/bin/zsh</string>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/local/bin:/opt/homebrew/bin:/Users/<<REPLACE_USER>>/.local/bin</string>
        <!-- Phase 4: feed_registry の notion-medium-db 巡回（discover-notion-medium）に必須。
             実際の API キー値に置換する（plist は ~/Library/LaunchAgents/ 配下＝gitignore 対象外に注意）。 -->
        <key>NOTION_API_KEY</key>
        <string><<REPLACE_NOTION_API_KEY>></string>
    </dict>

    <key>StandardOutPath</key>
    <string><<REPLACE_REPO_PATH>>/wiki-vault/.discover-watchlist-out.log</string>
    <key>StandardErrorPath</key>
    <string><<REPLACE_REPO_PATH>>/wiki-vault/.discover-watchlist-err.log</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
